# billometer/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
[program:{{ service_name }}]
directory={{ service.dir.base }}
environment=PATH="{{ service.dir.base }}/bin",PYTHONPATH='{{ service.dir.base }}/lib/python2.6/site-packages'
{%- if service.name == 'api' %}
command=/srv/billometer/bin/gunicorn_start
{%- endif %}
{%- if service.name == 'collector' %}
command={{ service.dir.base }}/bin/celery worker -s {{ service.dir.base }}/logs/celerybeat-schedule -B -E --hostname=collector@{{ grains.fqdn }} --loglevel=INFO --concurrency=1 --config=celery_config --workdir={{ service.dir.base }}/site --logfile={{ service.dir.base }}/logs/collector.log
{%- endif %}
{%- if service.name == 'monitor' %}
command=/srv/billometer/site/manage.py celerycam
{%- endif %}
stdout_logfile={{ service.dir.base }}/logs/{{ service.name }}_app.log 
stderr_logfile={{ service.dir.base }}/logs/{{ service.name }}_error.log
user=billometer
autostart=true
autorestart=true
