# billometer/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- from "billometer/map.jinja" import server with context %}

[program:{{ service_name }}]
directory={{ server.dir.base }}
{%- if service.name == 'api' %}
command={{ server.dir.base }}/bin/gunicorn_start.sh
{%- endif %}
{%- if service.name == 'collector' %}
command={{ server.dir.base }}/bin/celery worker -s /var/log/billometer/celerybeat-schedule -B -E --hostname=collector@{{ grains.fqdn }} --loglevel=INFO --concurrency=1 --config=celery_config --workdir={{ server.dir.base }}/bin --logfile=/var/log/billometer/collector.log
{%- endif %}
{%- if service.name == 'monitor' %}
command={{ server.dir.base }}/bin/manage.py celerycam
{%- endif %}
stdout_logfile=/var/log/billometer/{{ service.name }}_app.log 
stderr_logfile=/var/log/billometer/{{ service.name }}_app.log
user=billometer
autostart=true
autorestart=true
