# billometer/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
[program:{{ service_name }}]
directory=/srv/billometer/site
environment=PATH="/srv/billometer/bin",PYTHONPATH='/srv/billometer/lib/python2.6/site-packages'
{%- if service.name == 'api' %}
command=/srv/billometer/bin/gunicorn_start
{%- endif %}
{%- if service.name == 'collector' %}
command=/srv/billometer/bin/celery worker -s /srv/billometer/logs/celerybeat-schedule -B -E --hostname=collector@{{ grains.fqdn }} --loglevel=INFO --concurrency=1 --config=celery_config --workdir=/srv/billometer/site --logfile=/srv/billometer/logs/collector.log
{%- endif %}
{%- if service.name == 'monitor' %}
command=/srv/billometer/site/manage.py celerycam
{%- endif %}
stdout_logfile=/srv/billometer/logs/{{ service.name }}_app.log 
stderr_logfile=/srv/billometer/logs/{{ service.name }}_error.log
user=billometer
autostart=true
autorestart=true
